# deployment/supervisor.conf
; HR System — Queue Worker
[program:hr-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/hr-system/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/hr-system/storage/logs/queue-worker.log
stopwaitsecs=3600

; HR System — Laravel Reverb (WebSocket)
[program:hr-reverb]
process_name=%(program_name)s
command=php /var/www/hr-system/artisan reverb:start --host=0.0.0.0 --port=8080
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/hr-system/storage/logs/reverb.log
stopwaitsecs=10

; HR System — Laravel Scheduler
[program:hr-scheduler]
process_name=%(program_name)s
command=/bin/bash -c "while true; do php /var/www/hr-system/artisan schedule:run --no-interaction >> /dev/null 2>&1; sleep 60; done"
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/hr-system/storage/logs/scheduler.log
